# nerv/etc/systemd/system/sshd-keygen@.service.d/override.conf
[Unit]
ConditionFileNotEmpty=
ConditionFileNotEmpty=|!/var/etc/ssh/ssh_host_%i_key

[Service]
ExecStart=
ExecStart=bash /usr/nerv/sshd-keygen.sh %i

