# 1_ipxe_server/coreos_updater/coreos_updater.timer
[Unit]
Description=Run the CoreOS Updater to pull new images.

[Timer]
OnCalendar=

[Install]
WantedBy=timers.target
