# nerv/etc/systemd/system/slurmctld.service.d/override.conf
[Unit]
ConditionHost=misato*
After=systemd-tmpfiles-setup.service

[Service]
ExecStart=
ExecStart=/usr/sbin/slurmctld -D -s $SLURMCTLD_OPTIONS