# nerv/etc/systemd/system/boot.mount
[Unit]
Description="Mount boot partition"

[Mount]
What=LABEL=boot
Where=/boot
Type=auto
Options=defaults

[Install]
WantedBy=local-fs.target