# nerv/etc/systemd/system/boot-efi.mount
[Unit]
Description="Mount efi partition"
After=boot.mount

[Mount]
What=LABEL=esp
Where=/boot/efi
Type=auto

[Install]
WantedBy=local-fs.target