# nerv/etc/profile.d/20-nerv-lmod.sh
#!/usr/bin/bash

export MODULEPATH_ROOT="/shared/spack/modules/"
export MODULEPATH="/shared/spack/modules/Core"
