# nerv/etc/profile.d/90-nerv-modules.sh
#!/usr/bin/env sh

MODULEPATH_ROOT="/shared/spack/modules/linux-centos9-x86_64/Core"
MODULEPATH="/shared/spack/modules/linux-centos9-x86_64/Core:$MODULEPATH"
