# nerv/etc/profile.d/modules.sh
#!/usr/bin/env sh

MODULEPATH="/shared/spack/modules/linux-centos9-x86_64/Core:$MODULEPATH"
