# nerv/etc/systemd/system/slurmdbd.service.d/override.conf
[Unit]
ConditionPathExists=
ConditionHost=misato*

[Service]
ExecStartPre=+/usr/bin/bash /usr/nerv/slurmdbdkey.sh